# roles/fedora-web/kinoite/build/files/synckinoite.cron
MAILTO=web-members@fedoraproject.org
25 * * * * apache /usr/local/bin/lock-wrapper synckinoite /usr/local/bin/synckinoite
