# roles/koji_builder/files/osbuildapi-update.cron
* * * * * root /usr/local/bin/lock-wrapper osbuildapi "/usr/local/bin/osbuildapi-update.sh"
