# roles/mailman3/files/systemd/hyperkitty-hourly.service
[Unit]
Description=hyperkitty hourly maintenance

[Service]
Type=oneshot
User=mailman
ExecStart=mailman-web runjobs hourly
