# roles/mailman3/files/systemd/hyperkitty-quarter_hourly.service
[Unit]
Description=hyperkitty maintenance every 15 minutes

[Service]
Type=oneshot
User=mailman
ExecStart=mailman-web runjobs quarter_hourly
