# roles/mailman3/files/systemd/hyperkitty-yearly.service
[Unit]
Description=hyperkitty yearly maintenance

[Service]
Type=oneshot
User=mailman
ExecStart=mailman-web runjobs yearly
