# roles/mailman3/files/systemd/hyperkitty-monthly.service
[Unit]
Description=hyperkitty monthly maintenance

[Service]
Type=oneshot
User=mailman
ExecStart=mailman-web runjobs monthly
