# roles/koji_builder/templates/osbuildapi-update-nft.cron
{% if host in groups['osbuild'] %}
*/5 * * * * root /usr/local/bin/lock-wrapper osbuildapi "/usr/local/bin/osbuildapi-update-nft.sh" 2>&1 | /usr/local/bin/nag-once osbuildapi-update-nft.sh 1d 2>&1
{% else %}
23 17 * * * root /usr/local/bin/lock-wrapper osbuildapi "/usr/local/bin/osbuildapi-update-nft.sh" 2>&1 | /usr/local/bin/nag-once osbuildapi-update-nft.sh 1d 2>&1
{% endif %}
