# files/download/download-sync.cron
# run twice daily rsync of download. but lock it
MAILTO=root@fedoraproject.org
00 11,23 * * * root /usr/local/bin/lock-wrapper sync-up-downloads "/usr/local/bin/sync-up-downloads"

