# roles/fasjson/files/fasjson-aliases.cron
00 19 * * * root /usr/local/bin/lock-wrapper fasjson-aliases "/usr/local/bin/fasjson-aliases.py 2>&1"
