# roles/bodhi2/critpath/files/critpath.cron
#!/bin/sh

# generate/update critical path data for Bodhi

cd /etc/bodhi/critpath
/root/releng-git/scripts/critpath.py all
