# roles/copr/frontend/files/systemd/system/httpd.service.d/50-oom.conf
# Keep httpd master running when child is OOM killed, rhbz#1947475
[Service]
OOMPolicy = continue
