# roles/copr/frontend/files/systemd/system/httpd.service.d/50-restart.conf
# When httpd is OOM killed or otherwise failed, start it again
[Service]
Restart = on-failure
RestartSec = 5s
