# roles/web-data-analysis/files/condense-mirrorlogs.cron
MAILTO=root@fedoraproject.org,nils@redhat.com,asaleh@redhat.com
0 06  * * * root /usr/local/bin/condense-mirrorlogs.sh > /dev/null
* 09  5 * * root /usr/bin/sleep $(( $RANDOM \% 3600 )); /usr/local/bin/compress_logs.sh > /dev/null