# roles/web-data-analysis/files/countme-update.cron
MAILTO=root@fedoraproject.org,mattdm@fedoraproject.org,smooge@smoogespace.com
0 09  * * * countme /usr/local/bin/countme-update.sh
