# roles/web-data-analysis/files/countme-update.cron
MAILTO=root@fedoraproject.org,mattdm@fedoraproject.org,nils@redhat.com,jantill@redhat.com
0 09  * * * countme /usr/local/bin/countme-update.sh > /dev/null
0 09  * * * countme /usr/local/bin/countme-centos-update.sh > /dev/null
