# roles/bodhi2/backend/files/bodhi-automated-pushes.cron
MAILTO=sysadmin-releng-members@fedoraproject.org
14 00 * * * apache /usr/local/bin/bodhi-automated-pushes.py
