# roles/grobisplitter/files/cents8-split.cron
10 10  * * * root /usr/local/bin/cents8-split.sh
0 0 * * 1 root /usr/bin/find /mnt/fedora/app/fi-repo/centos/stream8-kojitarget -mindepth 1 -maxdepth 1 -type d -mtime +14 -exec rm -rf '{}' \;
