# roles/mailman3/templates/mailmanweb.service.j2
[Unit]
Description= Mailman Web ui
DefaultDependencies=no
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=mailman
Group=mailman
WorkingDirectory=/etc/mailman3
ExecStart=/usr/bin/gunicorn --max-requests 50000 --timeout 1000 --workers={{ mailman_gunicorn_workers }} --log-devel=debug --forwarded-allow-ips="*" mailman_web.wsgi
ExecStop=/bin/kill -s TERM $MAINPID
