# roles/copr/frontend/files/systemd/system/fm-consumer@copr_messaging.service.d/50-copr.conf
[Service]
# We still need to move Frontend code to default system paths.
Environment=PYTHONPATH=/usr/share/copr/coprs_frontend

# When the service is killed or failed, start it again
Restart = on-failure
RestartSec = 5s

# The default fedora messaging service file runs under root
User=copr-fe
Group=copr-fe
