# roles/copr/frontend/files/systemd/system/fm-consumer@copr_messaging.service.d/50-restart.conf
# When the service is killed or failed, start it again
[Service]
Restart = on-failure
RestartSec = 5s
